# infrastructure/systemd/morning-digest.timer
[Unit]
Description=Run morning email digest every day at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
